# etc/agent/init/hyperswitch-cleanup.conf
# vim:set ft=upstart ts=2 et:
description "hyperswitch clean-up"
author "Lionel Zerbib <lionel.zerbib@huawei.com>"

start on starting neutron-ovs-cleanup

pre-start script
    echo "hyperswitch-cleanup pre-start"
    mkdir -p /var/log/hypernet
    . /opt/hypernet/bin/activate
    /opt/hypernet/bin/hyperswitch-cleanup --log-file /var/log/hypernet/hyperswitch-cleanup.log --config-file /etc/hyperswitch/hyperswitch.conf
end script
